# rootfs/etc/services.d/dns-hack/run
#!/usr/bin/env sh
if [ ! -f /etc/dnsmasq-resolv.conf ]
then
    cp -f /etc/resolv.conf /etc/dnsmasq-resolv.conf
    echo "nameserver 127.0.0.1" > /etc/resolv.conf
fi
if [ ! -f /etc/hosts.orig ]
then
    cp /etc/hosts /etc/hosts.orig
fi

env_vars=$(env | grep ".*_NAME=" | cut -d= -f1 | tr '\n' ' ')

cp -f /etc/hosts.orig /tmp/hosts
for env_var in $env_vars
do
  link=${env_var%_NAME}
  domain=$(cat /etc/dnsmasq-resolv.conf | grep search | cut -d' ' -f2)
  nameserver=$(cat /etc/dnsmasq-resolv.conf | grep nameserver | head -1 | cut -d' ' -f2)
  ip=$(nslookup "${link}.${domain}" ${nameserver}  | grep Address | tail -1 | cut -d: -f2  | cut -d' ' -f2)
  if [ -n "$ip" ]
  then
    echo "${ip} ${link}.${domain}" >> /tmp/hosts
  else
     logger "ip ${link}.${domain} skipped, it didn't resolve."
  fi
done

cp -f /tmp/hosts /etc/hosts.links

sleep 2147483647