# rootfs/etc/services.d/dns-hack/run
#!/usr/bin/with-contenv ash

touch /etc/dnsmasq-resolv.conf

echo "DNS : Waiting for dnsmasq"
while ! ps -ef | grep -v grep | grep dnsmasq-resolv.conf  &> /dev/null
do
    sleep 1
done

echo "DNS : Starting monitoring loop"
while true
do
    if  ( cat /etc/resolv.conf | grep "nameserver 127.0.0.1" &> /dev/null )
    then
        sleep 30
    else
        cp -f /etc/resolv.conf /etc/dnsmasq-resolv.conf
        echo "nameserver 127.0.0.1" > /etc/resolv.conf
        echo "nameserver 8.8.8.8" >> /etc/resolv.conf
    fi

    cp -f /etc/hosts.orig /tmp/hosts
    if [ -n "$TUTUM_CONTAINER_FQDN" ]
    then
        . /bin/get_hosts_from_tutum.sh
        . /bin/tutum_dns_hack.sh
    else
        . /bin/non_tutum_dns_hack.sh
    fi

     cp -f /etc/hosts /tmp/hosts.orig

     sort -u < /tmp/hosts > /etc/hosts

     if ! diff /etc/hosts /tmp/hosts.orig &> /dev/null
     then
        echo "Updating DNS values"
        killall -HUP dnsmasq  &> /dev/null
     fi
     sleep 5
done
