# rootfs/etc/services.d/dns-hack/run
#!/usr/bin/with-contenv sh

touch /etc/dnsmasq-resolv.conf

echo "Waiting for dnsmasq"
while ! ps -ef | grep -v grep | grep dnsmasq-resolv.conf
do
    sleep 1
done

echo "Starting DNS monitoring loop"
while true
do
    if  ( cat /etc/resolv.conf | grep "nameserver 127.0.0.1" )
    then
        sleep 30
        continue
    else
        cp -f /etc/resolv.conf /etc/dnsmasq-resolv.conf
        echo "nameserver 127.0.0.1" > /etc/resolv.conf
    fi
    . /bin/dns_hack.sh

     if ! diff /tmp/hosts /etc/hosts.links
     then
        echo "Updating DNS values"
        cp -f /tmp/hosts /etc/hosts.links
        killall -HUP dnsmasq
     else
        echo "No change to DNS values"
     fi

    sleep 30
done