# rootfs/etc/services.d/dns/run
#!/usr/bin/with-contenv sh
exec go-dnsmasq --default-resolver --append-domain --hostsfile=/etc/hosts ${DNSMASQ_ARGS} 2>&1 | logger
