# rootfs/etc/services.d/dns/run
#!/bin/sh -x
exec dnsmasq ${DNSMASQ_ARGS:- --server=/node.tutum.io/162.159.24.22 --server=/cont.tutum.io/162.159.24.22 --server=/svc.tutum.io/162.159.25.72} --no-daemon 2>&1 | logger