# rootfs/etc/services.d/dns/run
#!/usr/bin/env sh
touch /etc/hosts.links
exec dnsmasq ${DNSMASQ_ARGS} --expand-hosts  --resolv-file=/etc/dnsmasq-resolv.conf --addn-hosts=/etc/hosts.links --no-daemon 2>&1 | logger
